# solutions/etc/systemd/system/getty@tty7.service.d/override.conf
# /etc/systemd/system/getty@tty1.service.d/override.conf
# systemctl --system enable getty@tty7.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin sh --noclear %I $TERM
